# new-api.service
# File path: /etc/systemd/system/frog-api.service
# sudo systemctl daemon-reload
# sudo systemctl start frog-api
# sudo systemctl enable frog-api
# sudo systemctl status frog-api
[Unit]
Description=Frog API Service
After=network.target

[Service]
User=ubuntu  # 注意修改用户名
WorkingDirectory=/path/to/frog-api  # 注意修改路径
ExecStart=/path/to/frog-api/frog-api --port 3000 --log-dir /path/to/frog-api/logs  # 注意修改路径和端口号
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
